# private_dot_config/autostart/picom.desktop
[Desktop Entry]
Exec=picom --experimental-backends
Name=picom
Type=Application
Version=1.0
