# private_dot_config/autostart/start wallpaper.desktop
[Desktop Entry]
Exec=/home/takase/.config/manage_wp start
Name=start wallpaper
Type=Application
Version=1.0
