# private_dot_config/autostart/urxvtd.desktop
[Desktop Entry]
Exec=urxvtd -q
Name=urxvtd
Type=Application
Version=1.0
