# private_dot_config/autostart/XAVA.desktop
[Desktop Entry]
Exec=xava
Name=XAVA
Type=Application
Version=1.0
